# logic/crontabs
# For more information see the manual pages of crontab(5) and cron(8)
# 
# m h  dom mon dow   command

*/5 * * * * /usr/bin/php /var/www/site/script.php >/dev/null 2>&1
